# knowledgehub-updater.service
[Unit]
Description=KnowledgeHub auto-updater

[Service]
Type=oneshot
ExecStart=/opt/knowledgehub/knowledgehub-updater.sh
